# src/bhunter.service
[Unit]
Description = XXX
 
[Service]
ExecStart=/usr/bin/contrSSH start
Type=fork
 
[Install]
WantedBy=multi-user.target